# .config/systemd/user/feh.service
[Unit]
Description=Feh service

[Service]
Type=oneshot
ExecStart=/usr/bin/bash /home/adi/.fehbg

